# scripts/yieldmap-prod.service
[Unit]
Description=YieldMap Production App
Documentation=https://github.com/your-repo/yieldmap
After=network.target
Wants=yieldmap-worker-prod.service

[Service]
Type=simple
User=cosmos
WorkingDirectory=/opt/yieldmap
Environment=NODE_ENV=production
EnvironmentFile=/opt/yieldmap/.env.production
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=10
KillMode=process

# Output to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=yieldmap-prod

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/yieldmap

[Install]
WantedBy=multi-user.target
